# rium-dosimeter.service
[Unit]
Description=Rium GM Dosimeter Reader - ASNR Fixed Station
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/Xx_Custom_Rium_GM_xX
ExecStart=/usr/bin/python3 /home/pi/Xx_Custom_Rium_GM_xX/read_dosimeter.py --send-data --production

# Restart policy
Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rium-dosimeter

[Install]
WantedBy=multi-user.target
